# recipes-ti/beagleboard/beaglebone-capes/cape.service
[Unit]
Description=Beaglebone cape support

[Service]
ExecStart=/usr/bin/cape.sh

[Install]
WantedBy=basic.target
