# recipes-ti/beagleboard/gadget-init/99-hokey-pokey.rules
ACTION=="remove", DEVPATH=="/devices/platform/omap/ti81xx-usbss/musb-hdrc.1/usb1/1-1",  RUN+="/usr/bin/hokey-pokey.sh"

